# planttree-iot.service
[Unit]
Description=Plant Tree IoT Server
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/planttree-iot
ExecStart=/usr/bin/dotnet /var/www/planttree-iot/PlantTreeIoTServer.dll
Restart=always
RestartSec=10
Environment=ASPNETCORE_ENVIRONMENT=Production
Environment=ASPNETCORE_URLS=http://0.0.0.0:5000

[Install]
WantedBy=multi-user.target